# supervisord.conf
[supervisord]
nodaemon=true

[program:ws_python]
command=python server_v2.py
stdout_logfile=./ws.log
stderr_logfile=./ws.log